# alpn312-x11/rootfs/etc/init.d/fluxbox
#!/sbin/openrc-run

# Launch fluxbox window manager.

command="/usr/bin/fluxbox"
command_background="yes"
pidfile="/run/fluxbox.pid"

depend() {
  need xvfb
}
