# alpn312-x11/rootfs/etc/init.d/x11vnc
#!/sbin/openrc-run

# Launch a VNC server used by the developer for debugging purpose.
# Note: this server has no password.

command="/usr/bin/x11vnc"
command_args="-display :1 -nopw -xkb -bg"

depend() {
  need net fluxbox
}
