# se-grid-standalone/rootfs/etc/init.d/selenium
#!/sbin/openrc-run

# Launch the Selenium standalone server

command="java"
command_args="-jar /opt/selenium.jar standalone"
command_background="yes"
pidfile="/run/selenium.pid"
